# docker/crontab
20 6,18 * * * cd /buildzone && ./build-mesa 2>&1 | tee -a log/mesa-"$(date +\%Y\%m\%d\%H\%M)" log/cron
